# etc/cron.d/coderdojochi-phonehome
# /etc/cron.d/coderdojochi-phonehome: crontab entry for coderdojochi-phonehome package

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/10 * * * * root /bin/sh /etc/init.d/coderdojochi-phonehome
