# etc/cron.d/coderdojochi-phonehome
# /etc/cron.d/coderdojochi-phonehome: crontab entry for coderdojochi-phonehome package

*/2 * * * * root /etc/init.d/coderdojochi-phonehome
