# etc/init/coderdojochi-phonehome.conf
console none

start on (local-filesystems and net-device-up IFACE!=lo)
stop on [!12345]

script
exec /etc/init.d/coderdojochi-phonehome
end script
